# root/etc/services.d/supervisord/run
#!/usr/bin/with-contenv bash

exec supervisord --nodaemon --configuration /config/supervisord/supervisord.conf